# iscsi/aaa.service
[Unit]
Description=aa service

[Service]
Type=forking
ExecStart=/root/run.sh
ExecStop=/root/stop.sh

[Install]
WantedBy=multi-user.target
